# modules/onedrive/dotfiles/systemd/user/noctalia.service.d/override.conf
[Service]
Environment="NOCTALIA_DEBUG=1"
